# openvpn/start-openvpn.service
[Service]
ExecStart=/home/pi/ovpn/start_vpn.sh
[Install]
WantedBy=network.target